# roles/services/templates/sync-iscsi-luns.service.j2
# Managed by Ansible — Pacemaker-managed oneshot service
# Restores client-facing LIO config from shared ZFS storage after pool import.
#
# This service is NOT enabled at boot — Pacemaker starts it as a
# cluster resource ordered after zfs-pool and before vip-iscsi.
# ExecStop saves the final client config to shared storage during planned failover.
[Unit]
Description=Sync ZFS zvols to iSCSI LUN mappings
After={{ _iscsi_target_service }}.service
Requires={{ _iscsi_target_service }}.service

[Service]
Type=oneshot
ExecStart=/root/iscsi-config-manager.py sync
ExecStop=/root/iscsi-config-manager.py save
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
